# sl-map-web/contrib/sl-map-web-backup.timer
[Unit]
Description=Daily SQLite backup for sl-map-web

[Timer]
OnCalendar=daily
RandomizedDelaySec=30m
Persistent=true
Unit=sl-map-web-backup.service

[Install]
WantedBy=timers.target
